# searxng/mcp-searxng.service
[Unit]
Description=SearXNG MCP Server (Streamable HTTP)
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=herzig_group
Environment=SEARXNG_URL=http://localhost:8080
Environment=MCP_TRANSPORT=streamable-http
Environment=MCP_HOST=0.0.0.0
Environment=MCP_PORT=8001
ExecStart=/usr/bin/python3 /home/herzig_group/searxng/mcp_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
